# terraform/modules/infrastructure/config/crontab/root.crontab
@reboot /usr/bin/vncserver :1
@reboot systemctl nessusd restart
@reboot service postgresql restart